# deploy/veron.service
[Unit]
Description=VERON 3270 Terminal Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=veron
Group=veron
WorkingDirectory=/opt/veron
EnvironmentFile=/etc/veron/env
ExecStart=/usr/bin/sbcl --non-interactive \
    --load /home/veron/quicklisp/setup.lisp \
    --load load.lisp \
    --eval '(veron:start :host "0.0.0.0")'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=veron

[Install]
WantedBy=multi-user.target
